# supervisord/supervisord-keystone.conf
[program:keystone]
command=/bin/bash -c "exec start-stop-daemon --start --chuid keystone --chdir /var/lib/keystone --name keystone --exec /usr/bin/keystone-all"